# rootfs/service/irssi/run
#!/bin/bash

# Start irssi client
export TERM=xterm
exec irssi -d
