# rootfs/service/irssi/run
#!/bin/bash

# Start irssi client
exec irssi
